# playbooks/templates/supervisor.ini
[program:{{ app_name }}]
command={{webapp_dir}}/{{app_name}}/venv/bin/gunicorn -w4 helloworld -b unix:/tmp/gunicorn.sock
directory={{ webapp_dir }}/{{ app_name }}/src
autostart=true
autorestart=true
